# cryp.service
[Unit]
Description=Cryp AI Assistant
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=__VENV_PYTHON__ __CRYP_DIR__/main.py
Restart=on-failure
RestartSec=5
WorkingDirectory=__CRYP_DIR__
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/%U
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
